# currency/distribution/systemd.service
[Unit]
Description=Odigos Demo Currency – Nginx vhost
After=network.target php8.3-fpm.service nginx.service
Requires=php8.3-fpm.service nginx.service

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
