# currency/distribution/systemd.service
[Unit]
Description=PHP-FPM and Nginx Service
After=network.target

[Service]
WorkingDirectory=/opt/odigos-demo-currency
ExecStart=/bin/sh -c "php-fpm & nginx -g 'daemon off;'"
KillMode=mixed
Restart=on-failure
RestartSec=1s
RestartSteps=1s
RestartMaxDelaySec=1min
Type=simple
User=odigos
Group=odigos
Environment="PORT=8085"

[Install]
WantedBy=multi-user.target
