# currency/distribution/systemd.service
[Unit]
Description=Odigos Demo Currency Service
After=network.target

[Service]
ExecStart=/usr/bin/php -S localhost:9000 -t /opt/odigos-demo-currency
KillMode=mixed
Restart=on-failure
RestartSec=1s
RestartSteps=1s
RestartMaxDelaySec=1min
Type=simple
User=odigos
Group=odigos
WorkingDirectory=/opt/odigos-demo-currency
Environment="PORT=9000"
Environment="GEOLOCATION_SERVICE_HOST=localhost:8086"

[Install]
WantedBy=multi-user.target
