# geolocation/distribution/systemd.service
[Unit]
Description=Odigos Demo Geolocation (Rails)
After=network.target

[Service]
Type=simple
User=odigos
Group=odigos
WorkingDirectory=/opt/odigos-demo-geolocation
Environment="BUNDLE_GEMFILE=/opt/odigos-demo-geolocation/Gemfile"
Environment="BUNDLE_PATH=/opt/odigos-demo-geolocation/vendor/bundle"
Environment="LD_LIBRARY_PATH=/opt/odigos-demo-geolocation/embedded-ruby/lib"
Environment="BUNDLE_WITHOUT=development:test"
Environment="BUNDLE_DEPLOYMENT=1"
Environment="BUNDLE_DISABLE_SHARED_GEMS=1"
Environment="RAILS_ENV=production"
Environment="PORT=8086"

ExecStart=/usr/bin/env bash -lc '\
  cd /opt/odigos-demo-geolocation && \
  /opt/odigos-demo-geolocation/embedded-ruby/bin/ruby \
    /opt/odigos-demo-geolocation/embedded-ruby/bin/bundle exec bin/rails db:prepare && \
  exec /opt/odigos-demo-geolocation/embedded-ruby/bin/ruby \
    /opt/odigos-demo-geolocation/embedded-ruby/bin/bundle exec bin/rails server --port=${PORT}'

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
